# gandalf-backup-cleanup.service
[Unit]
Description=Gandalf backup cleanup (local purge + cloud retention)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/gandalf/backup.env
ExecStart=/usr/bin/flock -n /run/gandalf-backup-cleanup.lock /usr/local/sbin/gandalf-backup-cleanup
TimeoutStartSec=10m
StandardOutput=journal
StandardError=journal
