# systemd/systemd/gandalf-backup-full.service
[Unit]
Description=Gandalf full backup (ALT A)
After=network-online.target
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/gandalf-backup full
