# systemd/systemd/gandalf-backup-inc.service
[Unit]
Description=Gandalf incremental backup (ALT A)
After=network-online.target
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/gandalf-backup inc
