# systemd/systemd/gandalf-backup-full.timer
[Unit]
Description=Run full Gandalf backup weekly
[Timer]
OnCalendar=Sun *-*-* 05:30
Persistent=true
[Install]
WantedBy=timers.target
