# systemd/systemd/gandalf-backup-cleanup.service
[Unit]
Description=Gandalf backup cleanup (retention + tmp)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/gandalf-backup-cleanup
