# systemd/dropins/gandalf-backup-inc.service.d/50-lock.conf
[Service]
ExecStart=
ExecStart=/usr/bin/flock -n /run/lock/gandalf-backup.lock /usr/local/sbin/gandalf-backup inc
