# systemd/dropins/gandalf-backup-inc.service.d/99-post.conf
[Service]
ExecStartPost=
ExecStartPost=/usr/local/sbin/gandalf-backup-route inc
