# systemd/dropins/gandalf-backup-cleanup.service.d/override.conf
[Service]
SyslogIdentifier=gandalf-backup-cleanup
StandardOutput=journal
StandardError=journal
